# einkdisplay-web.service
[Unit]
Description=E-Ink Display Web UI
After=network.target

[Service]
User=<user>
WorkingDirectory=/home/<user>/einkdisplay
ExecStart=/home/<user>/einkdisplay-venv/bin/python web/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
